# accounting/caso/caso.cron
# Run cASO every hour 
14 * * * * root /usr/local/bin/caso-extract.sh >> /var/log/caso.log 2>&1
