# bdii/bdii.ini
[supervisord]
nodaemon = true
loglevel = debug

[program:bdii]
command = pidproxy /var/run/bdii/bdii-update.pid /usr/local/bin/bdii-update
autostart = true
autorestart = true
user = ldap

[program:ldap]
command = /usr/sbin/slapd -f /etc/bdii/bdii-slapd.conf  -h ldap://0.0.0.0:2170 -u ldap -d 0
autostart = true
autorestart = true
