# conf/cron.d/atrope
# Run atrope every 12 hours
23 */12 * * * /usr/local/bin/atrope-dispatch.sh >> /var/log/atrope.log 2>&1
