# conf/cron.d/caso
# Run cASO every hour 
14 * * * * /usr/local/bin/caso-extract.sh
