# conf/cron.d/caso
# Run cASO every hour 
14 * * * * docker run -v /etc/voms.json:/etc/voms.json -v /etc/caso:/etc/caso -v /var/spool/caso:/var/spool/caso -v /var/spool/ssm:/var/spool/ssm egifedcloud/caso
