# conf/cron.d/ssmsend
# Send SSM records every 6 hours
30 */6 * * * docker run -v /etc/grid-security:/etc/grid-security -v /var/spool/apel:/var/spool/apel ssm ssmsend

