# conf/cron.d/ssmsend
# Send SSM records every 6 hours
30 */6 * * * /usr/local/bin/ssm-send.sh >> /var/log/ssm.log 2>&1

